# deploy/afficho.service
[Unit]
Description=Afficho Digital Signage Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=afficho
Group=afficho

# Apply staged auto-update (if any) before starting the daemon.
# Runs as root (before User= takes effect) so it can replace the binary.
ExecStartPre=/bin/sh -c '\
  if [ -f /var/lib/afficho/.update/pending ]; then \
    mv /var/lib/afficho/.update/afficho.new /usr/local/bin/afficho && \
    rm -f /var/lib/afficho/.update/pending && \
    echo "Applied staged update"; \
  fi'

ExecStart=/usr/local/bin/afficho -config /etc/afficho/config.toml
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/afficho /var/log/afficho

[Install]
WantedBy=multi-user.target
